# relay-control.service
[Unit]
Description=8-Relay Control Web Service
After=network.target

[Service]
Type=simple
User=tech
Group=gpio
WorkingDirectory=/home/tech/8-relay
# The ExecStart line now uses gunicorn for a more robust production server
ExecStart=/home/tech/8-relay/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:5000 app:app
Restart=always
RestartSec=10

# Security settings from the more robust service file
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Grant write access to necessary directories
ReadWritePaths=/var/log/relay_control /home/tech/8-relay

# Resource limits from the more robust service file
CPUQuota=50%
MemoryLimit=256M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=relay-control

# Audio permissions
SupplementaryGroups=audio

[Install]
WantedBy=multi-user.target
